# ss_firewalld.ini
[program:ssr_firewalld]
command=/root/ssmgr/ss_iptables_service
autorestart=true
autostart=true
user=root
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /root/ssmgr/log/ssmgr.log
